# crontab
# m 	h  	dom 	mon 	dow   	command
*/15	*	*	*	*	exec /var/www/colour.petersobot.com/run 
